# root/etc/services.d/mount/run
#!/bin/bash

umask 022
echo "Started mounting..."
exec mount